# raspberry-pi/services/chromium-kiosk-props.service
[Unit]
Description=Chromium Kiosk Mode - Props Pi (World Map Display)
After=graphical.target puzzle-2-world-map.service
Wants=graphical.target

[Service]
Type=simple
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/home/pi/hollywoodHQ/raspberry-pi/scripts/start-kiosk-props.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
