# raspberry-pi/services/puzzle-4-vehicle.service
[Unit]
Description=Mission Hollywood - Puzzle 4 (Vehicle Selector)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/hollywoodHQ/puzzle-4-vehicle
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
